# deploy/faratech-api.service
[Unit]
Description=Faratech News API (FastAPI)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/faratech-news-backend
EnvironmentFile=/opt/faratech-news-backend/.env
ExecStart=/opt/faratech-news-backend/.venv/bin/uvicorn main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
